# charm/launchpad-librarian/templates/launchpad-librarian.service.j2
# This service is really a systemd target, but we use a service since
# targets cannot be reloaded.  See launchpad-librarian@.service for instance
# configuration.

[Unit]
Description=Launchpad librarian

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
ExecReload=/bin/true

[Install]
WantedBy=multi-user.target

